# docker/dockerfile/mongodb/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:mongo]
command=/opt/mongodb-3.4.3/bin/mongod -f /mongodb/mongodb.conf