# docker/dockerfile/activemq/data/supervisord.conf
[supervisord]
nodaemon=true

[program:activemq]
user=root
command=/opt/activemq/bin/activemq start