# docker/dockerfile/redis/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:redis]
command=/opt/redis-3.2.8/src/redis-server /redis/redis.conf