# docker/mongodb/supervisord.conf
[supervisord]
nodaemon=true

#[program:sshd]
#command=/usr/sbin/sshd -D

[program:mongo]
command=/opt/mongodb-linux-x86_64-4.0.9/bin/mongod -f /opt/data/mongod.conf
