# docker/activemq/supervisord.conf
[supervisord]
nodaemon=true

[program:activemq]
user=root
command=/opt/apache-activemq-5.15.8/bin/activemq start